# services/screensaverd.service
[Unit]
Description=Idle screensaver watcher (launches YouTube saver on VT 7)
After=default.target

[Service]
Type=simple
ExecStart=/usr/local/bin/screensaver_watch.py
Environment=SCREENSAVER_IDLE=600 SCREENSAVER_REQUIRE_AUDIO_IDLE=0 SCREENSAVER_VT=7
Restart=always
RestartSec=2
# Allow reading /dev/input/event*
SupplementaryGroups=input

[Install]
WantedBy=default.target